# caps2esc.service
[Unit]
Description=Improve the behaviour of the CapsLock key

[Service]
ExecStart=XXXXXXXX
Type=exec

[Install]
WantedBy=default.target
